# resmon_scripts/service_units/com.resmon.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.resmon.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{PYTHON}}</string>
        <string>-m</string>
        <string>implementation_scripts.daemon</string>
        <string>{{PORT}}</string>
    </array>
    <key>WorkingDirectory</key>
    <string>{{WORKDIR}}</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>{{WORKDIR}}</string>
        <!-- Augment PATH so subprocess lookups (notably xelatex/pdflatex
             for routine-fired execution PDF bundles) can find binaries
             that the default launchd PATH does not include. -->
        <key>PATH</key>
        <string>/Library/TeX/texbin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>{{LOG_DIR}}/resmon-daemon.out.log</string>
    <key>StandardErrorPath</key>
    <string>{{LOG_DIR}}/resmon-daemon.err.log</string>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
